# tools/supervisor/cowork_event_watcher.conf
; Supervisord program block for cowork event watcher (slice cowork.004).
;
; Purpose:
;   Polls TDA signal journal + bias snapshot, writes
;   /opt/smc-cowork/triggers/event_flag.json on changes. Consumed by the
;   cowork bot's STEP 0a secondary trigger (cowork.runner.should_run_now).
;
; Install (VPS):
;   sudo cp tools/supervisor/cowork_event_watcher.conf \
;           /etc/supervisor/conf.d/cowork_event_watcher.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start cowork_event_watcher
;
; Pre-req: token with `read` scope must be present in
; COWORK_EVENT_WATCHER_TOKEN env var (rotated alongside other cowork tokens).
;
; Logs:
;   /var/log/smc-v3/cowork_event_watcher.stdout.log
;   /var/log/smc-v3/cowork_event_watcher.stderr.log
;
; Triggers dir must be writable by `ubuntu`:
;   sudo mkdir -p /opt/smc-cowork/triggers
;   sudo chown ubuntu:ubuntu /opt/smc-cowork/triggers

[program:cowork_event_watcher]
command=/opt/smc-v3/.venv/bin/python -m tools.cowork.event_watcher
directory=/opt/smc-v3
autostart=true
autorestart=true
startsecs=5
startretries=10
stopsignal=TERM
stopwaitsecs=15
user=ubuntu
stdout_logfile=/var/log/smc-v3/cowork_event_watcher.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/smc-v3/cowork_event_watcher.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=
    COWORK_TRIGGERS_DIR="/opt/smc-cowork/triggers",
    COWORK_SIGNALS_DIR="/opt/smc-v3/data_v3/_signals",
    COWORK_EVENT_WATCHER_API_BASE="http://127.0.0.1:8000",
    COWORK_EVENT_WATCHER_INTERVAL_S="30",
    COWORK_EVENT_WATCHER_SYMBOL="XAU/USD",
    COWORK_EVENT_WATCHER_LOG_LEVEL="INFO"
; COWORK_EVENT_WATCHER_TOKEN must be added separately (e.g. via
; /etc/supervisor/conf.d/cowork_event_watcher.local.conf or sourced env file).
