# tools/smc-v3.supervisor.conf
[program:smc-v3]
command=/opt/smc-v3/.venv/bin/python -m app.main --mode all
directory=/opt/smc-v3
user=ubuntu
autostart=false
autorestart=true
startretries=3
startsecs=15
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/var/log/smc-v3/stdout.log
stderr_logfile=/var/log/smc-v3/stderr.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=HOME="/home/ubuntu",PATH="/opt/smc-v3/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
