# tools/smc-v3.supervisor.conf
; =============================================================================
; SMC Platform v3 — Modular Supervisor Config
; =============================================================================
; 6 independent programs instead of 1 monolith.
; Benefits:
;   - `sudo supervisorctl restart smc:smc-ws` restarts ONLY ws_server
;   - No duplicate/orphan processes from Python 3.14 trampoline
;   - Each process can be stopped/started independently
;   - `sudo supervisorctl status` shows per-process health
;
; Commands:
;   sudo supervisorctl restart smc:smc-ws          # UI dev iteration
;   sudo supervisorctl restart smc:smc-fxcm        # broker reconnect
;   sudo supervisorctl restart smc:*               # full restart
;   sudo supervisorctl stop smc:smc-fxcm           # disable FXCM
;   sudo supervisorctl status                      # see all
; =============================================================================

; --- WS Server (UI backend, port 8000) ---
; Restart this when deploying UI changes or ws_server code changes.
[program:smc-ws]
command=/opt/smc-v3/.venv/bin/python -u -m runtime.ws.ws_server
directory=/opt/smc-v3
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/smc-v3/ws_server.stdout.log
stderr_logfile=/var/log/smc-v3/ws_server.stderr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
environment=HOME="/home/ubuntu",PYTHONUNBUFFERED="1",PATH="/opt/smc-v3/.venv/bin:/usr/local/bin:/usr/bin"

; --- FXCM Data Pipeline (broker_sidecar .venv37 + m1_ingestion) ---
; Stable, rarely needs restart.
[program:smc-fxcm]
command=/opt/smc-v3/.venv/bin/python -u -m app.main --mode m1_poller --stdio files --log-dir /var/log/smc-v3
directory=/opt/smc-v3
user=ubuntu
autostart=true
autorestart=true
startretries=3
startsecs=15
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/var/log/smc-v3/fxcm.stdout.log
stderr_logfile=/var/log/smc-v3/fxcm.stderr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
environment=HOME="/home/ubuntu",PYTHONUNBUFFERED="1",PATH="/opt/smc-v3/.venv/bin:/opt/smc-v3/.venv37/bin:/usr/local/bin:/usr/bin"

; --- FXCM Tick Stream ---
[program:smc-ticks]
command=/opt/smc-v3/.venv37/bin/python -u -m runtime.ingest.tick_publisher_fxcm
directory=/opt/smc-v3
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/smc-v3/ticks.stdout.log
stderr_logfile=/var/log/smc-v3/ticks.stderr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
environment=HOME="/home/ubuntu",PYTHONUNBUFFERED="1",PYTHONPATH="/opt/smc-v3",PATH="/opt/smc-v3/.venv37/bin:/usr/local/bin:/usr/bin"

; --- Tick Preview Worker ---
[program:smc-preview]
command=/opt/smc-v3/.venv/bin/python -u -m runtime.ingest.tick_preview_worker
directory=/opt/smc-v3
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/smc-v3/preview.stdout.log
stderr_logfile=/var/log/smc-v3/preview.stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
environment=HOME="/home/ubuntu",PYTHONUNBUFFERED="1",PATH="/opt/smc-v3/.venv/bin:/usr/local/bin:/usr/bin"

; --- Binance Ingest Worker ---
[program:smc-binance]
command=/opt/smc-v3/.venv/bin/python -u -m runtime.ingest.binance_ingest_worker
directory=/opt/smc-v3
user=ubuntu
autostart=true
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/smc-v3/binance.stdout.log
stderr_logfile=/var/log/smc-v3/binance.stderr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
environment=HOME="/home/ubuntu",PYTHONUNBUFFERED="1",PATH="/opt/smc-v3/.venv/bin:/usr/local/bin:/usr/bin"

; --- Binance Tick Publisher ---
[program:smc-binance-ticks]
command=/opt/smc-v3/.venv/bin/python -u -m runtime.ingest.binance_tick_publisher
directory=/opt/smc-v3
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/smc-v3/binance_ticks.stdout.log
stderr_logfile=/var/log/smc-v3/binance_ticks.stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
environment=HOME="/home/ubuntu",PYTHONUNBUFFERED="1",PATH="/opt/smc-v3/.venv/bin:/usr/local/bin:/usr/bin"

; --- Group: manage all with smc:* ---
[group:smc]
programs=smc-ws,smc-fxcm,smc-ticks,smc-preview,smc-binance,smc-binance-ticks
